# rustdesk/rootfs/etc/s6-overlay/s6-rc.d/hbbr/run
#!/command/with-contenv bashio
bashio::log.info "RUNNING HBBR"
cd /config
PARAMS="-k _"
/usr/bin/hbbr $PARAMS
