# rustdesk/rootfs/etc/s6-overlay/s6-rc.d/hbbs/run
#!/command/with-contenv bashio
bashio::log.info "RUNNING HBBS"
sleep 2
cd /config
PARAMS="-k _"
[ ! "$RELAY" = "null" ] && PARAMS="-r $RELAY $PARAMS"
/usr/bin/hbbs $PARAMS
